# scripts/ping_parser/crontab
*/10 * * * * bash PATH_TO_REPOS/scripts/ping_parser